# crontab
0 0 * * * cd /app && make clean-for-update all
